# systemd/cloudflare-r2-metrics.service
[Unit]
Description=Collect Cloudflare R2 metrics into InfluxDB v2
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=laurent
Group=laurent
WorkingDirectory=/home/laurent/script/cloudflare-r2
ExecStart=/usr/bin/env bash /home/laurent/script/cloudflare-r2/cloudflare-r2-metrics.sh
